# measurement/scaphandre/scaphandre.service
[Unit]
Description=Scaphandre Energy Metrics Exporter
After=network.target

[Service]
ExecStart=/usr/bin/scaphandre prometheus --port 8083
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
