# measurement/kepler/kepler.service
[Unit]
Description=Kepler Power Exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kepler --config.file /etc/kepler/kepler.yaml
Restart=always
User=root

[Install]
WantedBy=multi-user.target