# measurement/prometheus/prometheus.service
[Unit]
Description=Prometheus
After=network.target scaphandre.service

[Service]
User=prometheus
Group=prometheus
ExecStart=/usr/local/bin/prometheus \
  --config.file=/etc/prometheus/prometheus.yml \
  --storage.tsdb.path=/var/lib/prometheus \
  --storage.tsdb.retention.time=90d \
  --web.listen-address=:9099
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
